# lib/udev/rules.d/30-sound.rules
##

SUBSYSTEM=="sound", ACTION=="add|change", RUN+="set_volume"
