# etc/init/cache.conf
# drop cache periodically

description "drop cache"
author "nazarko"

start on (started mountall)
stop on (runlevel [!2345])

expect fork

exec /usr/bin/cachemon
