# etc/init/hama.conf
# hama
#

description	"hama mce"
author		"None"
start on (started xbmc and hama)
stop on (hamastop or runlevel [!2345])

task

exec /usr/sbin/hama_mce -d > /var/log/hama.log
