# etc/init/samba.conf
start on xbmc-started

script
	/etc/init.d/samba start
	exit 0
end script
